# target/linux/zynq/base-files/etc/board.d/gpio
#!/bin/sh /etc/rc.common

echo "Starting GPIO config" 
cd /sys/class/gpio
echo "937" > export
echo "938" > export
echo "out" > gpio937/direction
echo "out" > gpio938/direction


